# recipes-dca/execsnoop/files/dca-execsnoop.service
[Unit]
Description=Trace new processes via exec() syscalls

[Service]
Type=simple
ExecStart=/usr/share/bcc/tools/execsnoop
StandardOutput=file:/run/dca-execsnoop.log

[Install]
RequiredBy=basic.target
